# etc/systemd/system/getty@tty1.service.d/override.conf
[Service]
TTYVTDisallocate=no
ExecStart=
ExecStart=-/usr/bin/mingetty %I linux --autologin neg
