# etc/profile.d/Mathematica.sh
#!/bin/sh
export MATHEMATICA_HOME="/usr/local/Wolfram/Mathematica/10.0"
