# etc/systemd/system/ppp@.service.d/ppp.conf
[Service]
ExecStartPre=/usr/sbin/ip link set dev enp7s0 up
