# etc/systemd/system/auditd.service.d/override.conf
[Unit]
ConditionKernelCommandLine=
ConditionKernelCommandLine=audit=1
